# scripts/madeena-monitor.service
[Unit]
Description=Madeena Server Monitor
Documentation=https://github.com/Madeena-software/madeena-server-monitor
After=network.target

[Service]
Type=simple
User=madeena-monitor
Group=madeena-monitor
WorkingDirectory=/opt/madeena-monitor
EnvironmentFile=/opt/madeena-monitor/.env
ExecStart=/opt/madeena-monitor/monitor
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=madeena-monitor

# Security hardening
# NoNewPrivileges is intentionally omitted so that AmbientCapabilities can
# be inherited by child processes (smartctl needs CAP_SYS_RAWIO/CAP_SYS_ADMIN
# for SMART ioctl access; keeping it would silently block those capabilities).
CapabilityBoundingSet=CAP_SYS_RAWIO CAP_SYS_ADMIN
AmbientCapabilities=CAP_SYS_RAWIO CAP_SYS_ADMIN
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/madeena-monitor

[Install]
WantedBy=multi-user.target
